# backend/server/debian/scripts/crontab
*/5 * * * * cd /home/ubuntu/;/home/ubuntu/backup.sh >> /var/log/backups.log 2>&1